# opt/systemd/user/aether-desktop-recovery.service
[Unit]
Description=AetherOS Desktop Recovery Monitor
Documentation=man:aether-desktop-recovery(1)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/share/aetheros/scripts/aether-desktop-recovery monitor
Restart=on-failure
RestartSec=5

# Don't restart display manager
# Don't kill user applications
KillMode=process

[Install]
WantedBy=graphical-session.target
